# infra/systemd/impose.service
[Unit]
Description=Impose ton anonymat
After=network.target
StartLimitIntervalSec=120

[Service]
Type=simple
Restart=always
RestartSec=20
User=debian
WorkingDirectory=REMOTE_PATH/impose_ton_anonymat
EnvironmentFile=REMOTE_PATH/impose_ton_anonymat/infra/secrets
ExecStart=REMOTE_PATH/impose_ton_anonymat/run.sh

[Install]
WantedBy=multi-user.target
