# infra/systemd/impose.service
[Unit]
Description=Impose ton anonymat
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=debian
WorkingDirectory=/home/debian/bots/impose_ton_anonymat
ExecStart=/home/debian/bots/impose_ton_anonymat/run.sh

[Install]
WantedBy=multi-user.target
